# package-scripts/flapjack/postinst
#!/bin/bash
#
# Perform necessary flapjack setup steps after package is installed.
#

INSTALLER_DIR=`dirname $0`
DEST_DIR=/opt/flapjack
CONFIG_DIR=/etc/flapjack
FLAPJACK_GEM_DIR=`ls -dr ${DEST_DIR}/embedded/lib/ruby/gems/*/gems/flapjack-* | tail -1`
PROGNAME=$(basename $0)

groupadd flapjack
useradd --system --gid flapjack --no-create-home flapjack
if [ -f /etc/redhat-release ]; then
  INITDDIR="/etc/rc.d/init.d"
else
  INITDDIR="/etc/init.d"
fi

function error_exit
{
  echo "${PROGNAME}: ${1:-"Unknown Error"}" 1>&2
  exit 1
}

# Set up config + logging directories
mkdir -p /etc/flapjack /var/log/flapjack /var/run/flapjack /var/lib/flapjack/redis-flapjack
chown -R flapjack:flapjack /var/log/flapjack /var/run/flapjack /var/lib/flapjack

if [ -e ${FLAPJACK_GEM_DIR}/etc/flapjack_config.yaml.example ]; then
  FLAPJACK_CONFIG_VERSION="1"
elif [ -e ${FLAPJACK_GEM_DIR}/etc/flapjack_config.toml.example ]; then
  FLAPJACK_CONFIG_VERSION="2"
fi

# Copy in default configuration if none exists already
if [ ! -e /etc/flapjack/flapjack_config.yaml ] ; then
  if [ "${FLAPJACK_CONFIG_VERSION}" == "1" ]; then
    echo "Creating /etc/flapjack/flapjack_config.yaml"
    cp ${FLAPJACK_GEM_DIR}/etc/flapjack_config.yaml.example ${CONFIG_DIR}/flapjack_config.yaml
  fi
fi

if [ ! -e /etc/flapjack/flapjack_config.toml ] ; then
  if [ "${FLAPJACK_CONFIG_VERSION}" == "2" ]; then
    echo "Creating /etc/flapjack/flapjack_config.toml"
    cp ${FLAPJACK_GEM_DIR}/etc/flapjack_config.toml.example ${CONFIG_DIR}/flapjack_config.toml
  fi
fi

for version in 3 4 ; do
  if [ ! -e /usr/local/lib/flapjackfeeder${version}.o ] ; then
    echo "Creating /usr/local/lib/flapjackfeeder${version}.o for Nagios ${version}"
    cp ${FLAPJACK_GEM_DIR}/flapjackfeeder${version}.o /usr/local/lib/flapjackfeeder${version}.o
  fi
done

if [ ! -e /usr/local/lib/flapjackfeeder.o ] ; then
  echo "Symlinking the Nagios version 3 flapjackfeeder to /usr/local/lib/flapjackfeeder.o for backwards compatibility"
  ln -s /usr/local/lib/flapjackfeeder3.o /usr/local/lib/flapjackfeeder.o
fi

if [ ! -e ${INITDDIR}/redis-flapjack ] ; then
  echo "Creating ${INITDDIR}/redis-flapjack"
  if [ -f /etc/redhat-release ]; then
    cp ${DEST_DIR}/embedded/etc/init.d/redis-flapjack-rpm ${INITDDIR}/redis-flapjack
  else
    cp ${DEST_DIR}/embedded/etc/init.d/redis-flapjack-deb ${INITDDIR}/redis-flapjack
  fi
  chmod u+x ${INITDDIR}/redis-flapjack
fi

if [ ! -e ${INITDDIR}/flapjack ] ; then
  echo "Creating ${INITDDIR}/flapjack"
  if [ "${FLAPJACK_CONFIG_VERSION}" == "1" ]; then
    cp ${DEST_DIR}/embedded/etc/init.d/v1/flapjack ${INITDDIR}/flapjack
  elif [ -f /etc/redhat-release ]; then
    cp ${DEST_DIR}/embedded/etc/init.d/v2/rpm/flapjack ${INITDDIR}/flapjack
  else
    cp ${DEST_DIR}/embedded/etc/init.d/v2/deb/flapjack ${INITDDIR}/flapjack
  fi
  chmod u+x ${INITDDIR}/flapjack
fi

if [ ! -e ${INITDDIR}/flapjack-nagios-receiver ] ; then
  echo "Creating ${INITDDIR}/flapjack-nagios-receiver"
  if [ "${FLAPJACK_CONFIG_VERSION}" == "1" ]; then
    cp ${DEST_DIR}/embedded/etc/init.d/v1/flapjack-nagios-receiver ${INITDDIR}/flapjack-nagios-receiver
  elif [ -f /etc/redhat-release ]; then
    cp ${DEST_DIR}/embedded/etc/init.d/v2/rpm/flapjack-nagios-receiver ${INITDDIR}/flapjack-nagios-receiver
  else
    cp ${DEST_DIR}/embedded/etc/init.d/v2/deb/flapjack-nagios-receiver ${INITDDIR}/flapjack-nagios-receiver
  fi
  chmod u+x ${INITDDIR}/flapjack-nagios-receiver
fi

if [ ! -e ${INITDDIR}/flapper ] ; then
  echo "Creating ${INITDDIR}/flapper"
  if [ "${FLAPJACK_CONFIG_VERSION}" == "1" ]; then
    cp ${DEST_DIR}/embedded/etc/init.d/v1/flapper ${INITDDIR}/flapper
  elif [ -f /etc/redhat-release ]; then
    cp ${DEST_DIR}/embedded/etc/init.d/v2/rpm/flapper ${INITDDIR}/flapper
  else
    cp ${DEST_DIR}/embedded/etc/init.d/v2/deb/flapper ${INITDDIR}/flapper
  fi
  chmod u+x ${INITDDIR}/flapper
fi

echo "Setting paths"
if [ ! -e /etc/profile.d/flapjack.sh ] ; then
  cat > /etc/profile.d/flapjack.sh <<'EOF'
PATH=$PATH:/opt/flapjack/bin
EOF
  chmod 644 /etc/profile.d/flapjack.sh
fi

if [ -f /etc/redhat-release ]; then
  echo "In order to have flapjack start on reboot, run the following commands:"
  echo "  chkconfig --add redis-flapjack"
  echo "  chkconfig --add flapjack"
  echo
  echo "You can now start redis and flapjack as follows:"
  echo "  ${INITDDIR}/redis-flapjack start"
  echo "  ${INITDDIR}/flapjack start"
else
  update-rc.d redis-flapjack defaults 94 06 >/dev/null || exit $?
  update-rc.d flapjack defaults 95 05 >/dev/null || exit $?
  echo "Starting redis-flapjack"
  ${INITDDIR}/redis-flapjack start || exit $?

  echo "Starting flapjack"
  ${INITDDIR}/flapjack start || exit $?
fi

echo "Thank you for installing flapjack!"
echo
echo "Log in again to have /opt/flapjack/bin added to your PATH or run:"
echo "  export PATH=$PATH:/opt/flapjack/bin"
echo

if [ -f /etc/redhat-release ]; then
  echo "Once you've started Redis and Flapjack (see above) you'll be able to visit http://localhost:3080/"
else
  echo "Flapjack should now be running. Visit http://localhost:3080/"
fi

exit 0
